# files/etc/uci-defaults/99-default-lang
#!/bin/sh
uci set luci.main.lang=zh_Hans
uci commit luci
exit 0
